# util/crontab
30 02 * * * /var/www/ireceptor_gateway/util/scripts/delete_old_downloads.sh
